# assetshold.service
[Unit]
Description=Assets Portfolio Management Application
After=network.target

[Service]
Type=simple
User=yang_server
Group=yang_server
WorkingDirectory=/home/yang_server/assetshold
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
Environment=MARKET_ENABLE=1
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assetshold

[Install]
WantedBy=multi-user.target